# scripts/autostart/task-try-sync.desktop
[Desktop Entry]
Name=Taskwarrior Try Sync
Exec=/home/jangho/.setup/scripts/bin/task-try-sync
Terminal=false
