# task/autostart/task-sync.desktop
[Desktop Entry]
Name=Task sync
Exec=task sync
Terminal=false
Type=Application
